# etc/cron.d/apt-fetch
# don't run this too often or you may get dropped by the mirrors
*/30 * * * * root /usr/local/bin/apt-fetch
